# py/extra/prometheus-node-exporter-supervisor.conf
[program:prometheus_node_exporter]
command=/opt/prometheus_node_exporter/node_exporter
stdout_logfile=/var/log/prometheus_node_exporter.log
stdout_logfile_backups=1
redirect_stderr=true
